# ansible/roles/hardening/templates/wifi-watchdog.service.j2
[Unit]
Description=WiFi Connectivity Watchdog
After=network.target

[Service]
ExecStart=/usr/local/bin/wifi-watchdog.sh
Restart=always

[Install]
WantedBy=multi-user.target
